# cookbooks/scratchpads/templates/default/solr-undertow.service.erb
[Unit]
Description=Solr-Undertow

[Service]
Type=simple
Restart=always
ExecStart=<%= node['scratchpads']['solr-undertow']['templates']['bash_script']['path'] %>
User=<%= node['scratchpads']['solr-undertow']['user'] %>
Group=<%= node['scratchpads']['solr-undertow']['group'] %>

[Install]
WantedBy=multi-user.target