# cookbooks/scratchpads/templates/default/varnish.service.erb
[Unit]
Description=Varnish HTTP accelerator

[Service]
Type=forking
LimitNOFILE=<%= node['varnish']['nfiles'] %>
LimitMEMLOCK=<%= node['varnish']['memlock'] %>
EnvironmentFile=/etc/default/varnish
ExecStartPre=/usr/sbin/varnishd -C -f <%= node['varnish']['dir'] %>/<%= node['varnish']['vcl_conf'] %>
ExecStart=/usr/sbin/varnishd $DAEMON_OPTS
ExecReload=/usr/share/varnish/reload-vcl

[Install]
WantedBy=multi-user.target
